# config/twittrouter.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=98
STOP=98

start() {
    echo "starting twittrouter..."
    /usr/bin/twittrouter
}

stop() {
    echo "twittrouter stopped."
    ps |grep twittrouter |awk '{print $1}'|xargs kill -9
}
